# cmd/misc/freedom-speedtest.service
[Unit]
Description=freedom-speedtest

[Service]
User=nobody
ExecStart=/usr/bin/freedom-speedtest -o mongodb
