# cmd/misc/freedom-speedtest.timer
[Unit]
Description=freedom-speedtest

[Timer]
OnCalendar=hourly

[Install]
WantedBy=default.target
